# nsca.conf
[program:nsca]
command=/usr/sbin/nsca -c /etc/nsca.cfg -f --daemon
autostart=true
autorestart=true
startretries=1
startsecs=1
